# fs/etc/services.d/dns/run
#!/usr/bin/execlineb -P

s6-setuidgid nonroot

/usr/sbin/named -f -g -4 -c /config-default/named.conf